# hipdaspy/cron_5_min.cron
*/5 * * * * /root/hipida-spy/hipda-spiders/schedled_task.sh